# infra/ansible/playbooks/roles/install_app/templates/sessionid_scrapper.timer.j2
[Unit]
Description=Run sessionid_scraper every 12 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=12h
Unit=sessionid_scraper.service

[Install]
WantedBy=timers.target
